# ourhomeconnected/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the ourhomeconnected service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

bashio::log.info "Starting Our Home Connected backend..."

# Run the FastAPI app
exec /opt/venv/bin/python -m ohc_backend
